# docker-files/supervisord-extra.conf
[program:django-migrations]
command=python /app/manage.py migrate
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:django-collect-static]
command=python /app/manage.py collectstatic --noinput
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:django-fixtures]
command=python /app/manage.py loaddata /app/fixtures/data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
